# systemd/dahua-event-relay.service
[Unit]
Description=Dahua EventHttpUpload relay to Home Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=DAHUA_HA_WEBHOOK=http://192.168.1.11:8123/api/webhook/dahua_acs_events
ExecStart=/usr/local/bin/dahua_event_relay.py --host 0.0.0.0 --port 8818
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
